# deploy/systemd/aisw-bot.service
# ai-steward-wiki life bot — current single-user deploy unit (ADR-009 / ADR-010).
#
# Runs as `bgs` using that account's default ~/.claude (subscription auth, ADR-009)
# and ~/.local/bin (claude + uv). This is the CURRENT model: per-user UID isolation
# (D-038) is DEFERRED (ADR-010) for the trusted, allowlisted, single-user life
# deployment. See aisw-bot.service.d038-deferred for the future multi-tenant unit.
#
# Critical differences from the D-038 unit:
#   - No dedicated `aisw-bot` service user, no CAP_SETUID, no per-CLI scopes/slices.
#   - NO ProtectHome — the D-038 unit's `ProtectHome=tmpfs` would hide ~/.claude
#     and break subscription auth.
#   - PATH must include ~/.local/bin so resolve_binary("claude") finds the CLI.
#
# Install (VPS): copy to /etc/systemd/system/aisw-bot.service, then
#   sudo systemctl daemon-reload && sudo systemctl enable --now aisw-bot
# Logs: sudo journalctl -u aisw-bot -f

[Unit]
Description=ai-steward-wiki Telegram bot (life WIKI assistant)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bgs
Group=bgs
WorkingDirectory=/home/bgs/works/ai-steward-wiki
Environment=HOME=/home/bgs
Environment=PATH=/home/bgs/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=AISW_CODEX_HOME=/var/lib/ai-steward-wiki/codex
ExecStart=/home/bgs/works/ai-steward-wiki/.venv/bin/python -m ai_steward_wiki
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aisw-bot

[Install]
WantedBy=multi-user.target
